# debian/aiv.service
[Unit]
Description=AIV Application
After=network.target

[Service]
User=aiv
Group=aiv
Type=simple
ExecStart=/usr/bin/aiv
WorkingDirectory=/var/lib/aiv
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
