# DEBIAN/postinst
#!/bin/bash
set -e

# Create user if it doesn't exist
if ! id -u aiv > /dev/null 2>&1; then
    adduser --system --group --no-create-home aiv
fi

# Set proper ownership and permissions
chown -R aiv:aiv /var/lib/aiv /etc/aiv /var/logs/aiv
chmod 755 /var/lib/aiv /etc/aiv /var/logs/aiv

# Register and start the service if systemd is available
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable aiv.service || true
    systemctl start aiv.service || true
fi

exit 0
